# deploy/supervisor/aduan-watch.conf
; Supervisor program for the Aduan realtime push watcher.
;
; Copy to /etc/supervisor/conf.d/aduan-watch.conf on the API server, adjust the
; absolute paths + `user`, then:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start aduan-watch:*
;
; This daemon polls the shared `aduans` table every ~2s and pushes FCM for newly
; created complaints (web + mobile). It is SEPARATE from the queue worker, which
; still handles mobile-originated assignment/update events.

[program:aduan-watch]
process_name=%(program_name)s_%(process_num)02d
; --sleep is clamped to 1–3s by the command; --limit caps rows handled per poll.
command=php /var/www/ppa-api-new/artisan aduan:watch --sleep=2 --limit=200
directory=/var/www/ppa-api-new
user=www-data
numprocs=1
autostart=true
autorestart=true
startsecs=3
; Give the current 2s poll time to finish; the command traps SIGTERM and stops
; within ~1s, so restarts/deploys never drop or duplicate an event.
stopsignal=TERM
stopwaitsecs=15
redirect_stderr=true
stdout_logfile=/var/www/ppa-api-new/storage/logs/aduan-watch.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
